# supervisord.d/lotterybot_AutoSendMessage.ini
[program:lotterybot_AutoSendMessage]

process_name=%(program_name)s_%(process_num)02d

command=php /home/wwwroot/bot.test.com/think queue:work --daemon --queue AutoSendMessage --delay 10 --tries 1 --sleep 5 --memory 64

autostart=true

directory=/home/wwwroot/bot.test.com/

autorestart=true

startsecs=10

startretries=3

numprocs=3