# supervisord.d/lotterybot_AutoSendPrize.ini
[program:lotterybot_AutoSendPrize]

process_name=%(program_name)s_%(process_num)02d

command=php /home/wwwroot/bot.test.com/think queue:listen --queue AutoSendPrize --delay 10 --tries 3 --sleep 15 --memory 64

autostart=true

directory=/home/wwwroot/bot.test.com/

autorestart=true

startsecs=10

startretries=3

numprocs=1